# packaging/breadmill.service
[Unit]
Description=Breadmill semantic search indexer
Documentation=https://github.com/breadway/breadsearch
After=default.target

[Service]
Type=simple
# Required (not just a perf tweak) when backend = "rocm": without a valid
# cache dir, ONNX Runtime's MIGraphX EP reads an uninitialized cache path and
# crashes on the first query instead of just recompiling every restart. Inert
# for cpu/npu/cuda backends. See README.md#gpu-backend-notes.
Environment=ORT_MIGRAPHX_MODEL_CACHE_PATH=%h/.cache/breadsearch/migraphx-cache
ExecStart=%h/.cargo/bin/breadmill
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
