# packaging/breadmill.service
[Unit]
Description=Breadmill semantic search indexer
Documentation=https://github.com/breadway/breadsearch
After=default.target

[Service]
Type=simple
# Uncomment if built with --features rocm: persists MIGraphX's compiled-kernel
# cache across restarts (each new sequence length otherwise costs a ~60-120s
# recompile). See README.md#gpu-backend-notes.
#Environment=ORT_MIGRAPHX_MODEL_CACHE_PATH=%h/.cache/breadsearch/migraphx-cache
ExecStart=%h/.cargo/bin/breadmill
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
